# mac/whisper-server.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.openclaw.whisper-server</string>
  <key>ProgramArguments</key>
  <array>
    <!-- Update these paths to match your system -->
    <string>/opt/homebrew/bin/whisper-server</string>
    <string>-m</string>
    <string>/usr/local/share/whisper.cpp/ggml-medium.en.bin</string>
    <string>--port</string>
    <string>8080</string>
    <string>--host</string>
    <string>0.0.0.0</string>
    <string>-t</string>
    <string>8</string>
  </array>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/var/log/whisper-server/whisper-server.log</string>
  <key>StandardErrorPath</key>
  <string>/var/log/whisper-server/whisper-server.log</string>
  <key>WorkingDirectory</key>
  <string>/usr/local/share/whisper.cpp</string>
</dict>
</plist>

<!-- NOTE: install.sh generates this file dynamically with correct paths.
     If you edit this file manually, re-run install.sh to apply changes. -->