# vm/services/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway Service for Lucy
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/openclaw
EnvironmentFile=/opt/openclaw/openclaw.json
# Assuming OpenClaw is installed globally or via npm
ExecStart=/usr/bin/node /opt/openclaw/index.js
Restart=always
RestartSec=10
StandardOutput=append:/var/log/openclaw-gateway.log
StandardError=append:/var/log/openclaw-gateway.log

[Install]
WantedBy=multi-user.target