# vm/services/openclaw-gateway.service
[Unit]
Description=OpenClaw Gateway (Lucy Agent)
After=network-online.target voice-transcriber.service
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/openclaw
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables are loaded by OpenClaw itself from its config file.
# If OpenClaw supports an .env file, set EnvironmentFile here.

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/openclaw/logs
PrivateTmp=true

[Install]
WantedBy=multi-user.target
