# vm/services/voice-transcriber.service
[Unit]
Description=Discord Voice Transcriber Bot
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=openclaw
Group=openclaw
WorkingDirectory=/opt/voice-transcriber
EnvironmentFile=/opt/voice-transcriber/.env
ExecStart=/usr/bin/node index.js
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal

# Hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/voice-transcriber/temp_audio
PrivateTmp=true

[Install]
WantedBy=multi-user.target
